# rc.local
#!/bin/bash

# Start MagicPoint slideshow on boot on Display 0 in multi-user mode
# Author: Christopher Markieta

RAMDISK='/dev/shm/'

# Run in RAM to prevent damaging microSD card
cp -R /home/markieta/Documents/magicpoint/ "$RAMDISK"

export DISPLAY=:0

X -noreset -nocursor -ac $DISPLAY &

setterm -blank 0 -powersave off

xset s off
xset s noblank
xset -dpms

while true
do
    mgp -F 1,0,0 -d 5 "$RAMDISK"/magicpoint/ehl.mgp &>/dev/null
done &
